# prod/nginx-php-fpm/conf/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB      
logfile_backups=10        
loglevel=info            
pidfile=/tmp/supervisord.pid
nodaemon=false             
minfds=1024               
minprocs=200             
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock 

[program:php-fpm]
command=/usr/sbin/php-fpm -R -F
autostart=true
autorestart=true
priority=5

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
